# templates/activemq.conf.erb
start on startup
respawn
script
  if [ -f "/etc/default/activemq" ]; then
    set -a
    . /etc/default/activemq
  fi
  . /etc/default/activemq
  exec /usr/ActiveMQ/bin/activemq-admin start
end script