# templates/activemq.conf.erb
start on startup
respawn
script
  set -a
  . /etc/default/activemq
  exec /usr/ActiveMQ/bin/activemq-admin start
end script